# systemd/korauth.service
[Unit]
Description=korauth - OpenKor Identity & Access Management Service
Documentation=https://github.com/OpenKorProject/korauth
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=korauth
Group=korauth
WorkingDirectory=/opt/openkor/korauth

# Environment variables
EnvironmentFile=/opt/openkor/korauth/config/korauth.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/openkor/korauth

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Startup
ExecStart=/opt/openkor/korauth/bin/korauth
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=korauth

[Install]
WantedBy=multi-user.target
